# operations/docker/all-in-one/s6-rc.d/pyroscope/run
#!/bin/sh
exec 2>&1
exec /usr/bin/pyroscope
