# scripts/packages/pyroscope-server.service
[Unit]
Description=Open source continuous profiling software
Documentation="https://pyroscope.io/"
After=network-online.target

[Service]
User=pyroscope
Group=pyroscope
EnvironmentFile=-/etc/default/pyroscope
ExecStart=/usr/bin/pyroscope server
Restart=on-failure

[Install]
WantedBy=multi-user.target
