# operations/docker/all-in-one/s6-rc.d/grafana/run
#!/bin/sh
exec 2>&1
exec exec s6-setuidgid grafana /run.sh
